# .config/systemd/user/kmonad-restarter.service
[Unit]
Description=Restarts KMonad

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart kmonad.service
